# debian/sl-modem-daemon.prerm
#!/bin/sh
set -e
if [ -x "/etc/init.d/sl-modem-daemon" ]; then
	if [ -x /usr/sbin/invoke-rc.d ] ; then
		invoke-rc.d sl-modem-daemon stop || true
	else
		/etc/init.d/sl-modem-daemon stop || true
	fi
fi

case "$1" in
  remove)
    rm -rf /var/lib/slmodem
  ;;
esac
#DEBHELPER#
