# debian/sl-modem-daemon.postrm
#!/bin/sh
set -e

case "$1" in
    purge)
      update-rc.d -f sl-modem-daemon remove > /dev/null 2>&1
      # Remove log directory
      rm -rf /var/log/slmodemd
    ;;
esac

#DEBHELPER#

exit 0
