# debian/sl-modem-daemon.rules
KERNEL=="slusb*", NAME="slusb*", GROUP="dialout", RUN+="/etc/init.d/sl-modem-daemon"
